# tuned-watch-restarter.service
[Unit]
Description=tuned-watch restarter

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart tuned-watch.service
