# tuned-watch-restarter.path
[Path]
PathChanged=%h/.config/tuned-watch.conf

[Install]
WantedBy=tuned-watch-restarter.service

