# tuned-watch-restarter.path
[Unit]
Description=Watch changes to tuned-watch.conf

[Path]
PathChanged=%h/.config/tuned-watch.conf

[Install]
WantedBy=default.target

